# supervisor-sidekiq.conf
[program:sidekiq]
command=bundle exec sidekiq
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/sidekiq.log
stderr_logfile=/var/log/sidekiq.error.log
